# platform/systemd/lg-person-vis-refresh.timer
[Unit]
Description=Refresh person identity/visibility cache every 15 minutes

[Timer]
OnBootSec=4min
OnUnitActiveSec=15min
AccuracySec=1min
Unit=lg-person-vis-refresh.service

[Install]
WantedBy=timers.target
