# platform/systemd/lg-wp-cron.timer
[Unit]
Description=WP cron tick (server cron; wp-config sets DISABLE_WP_CRON)

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
